# home_scripts/scripts/jack-autostart.service
[Unit]
Description=Autostarts jack on pipewire

[Service]
Type=oneshoot
ExecStart="pw-jack jack_control start"

[Install]
WantedBy=user.target
